# parent/vsock-proxies.service
# systemd service unit for vsock proxies (Enclave -> External Services)
#
# These proxies allow the enclave to connect to services via vsock.
# The enclave connects to vsock ports which are forwarded to actual services.
#
# Traffic flow (PostgreSQL example):
#   Enclave localhost:5432 -> vsock CID 3 port 5432 -> This proxy -> localhost:5432
#
# Install:
#   sudo cp vsock-proxies.service /etc/systemd/system/
#   sudo systemctl daemon-reload
#   sudo systemctl enable vsock-proxies
#   sudo systemctl start vsock-proxies

[Unit]
Description=vsock Proxies for Nitro Enclave Egress
Documentation=https://github.com/aws/aws-nitro-enclaves-cli
After=network.target postgresql.service redis.service
Wants=postgresql.service

[Service]
Type=forking
User=root
Group=root

# Start all vsock proxies
ExecStart=/bin/bash -c '\
    /usr/bin/vsock-proxy 5432 localhost 5432 --config /etc/vsock-proxy/config.yaml & \
    /usr/bin/vsock-proxy 6379 localhost 6379 --config /etc/vsock-proxy/config.yaml & \
    /usr/bin/vsock-proxy 443 kms.us-east-1.amazonaws.com 443 --config /etc/vsock-proxy/config.yaml & \
    /usr/bin/vsock-proxy 8332 localhost 8332 --config /etc/vsock-proxy/config.yaml & \
    /usr/bin/vsock-proxy 8333 localhost 8333 --config /etc/vsock-proxy/config.yaml & \
    /usr/bin/vsock-proxy 8334 localhost 8334 --config /etc/vsock-proxy/config.yaml & \
    wait'

# Stop all proxies
ExecStop=/usr/bin/pkill -f vsock-proxy

# Restart on failure
Restart=always
RestartSec=5

# PID file for tracking
PIDFile=/run/vsock-proxies.pid

# Resource limits
LimitNOFILE=65535

[Install]
WantedBy=multi-user.target
