# parent/systemd/mint-creds-sender.service
[Unit]
Description=Mint Enclave credentials sender
After=mint-enclave.service
Requires=mint-enclave.service

[Service]
Type=oneshot
# Script retries internally, no need for ExecStartPre sleep
ExecStart=/opt/mint-enclave/send-credentials.sh
StandardOutput=journal
StandardError=journal
TimeoutStartSec=120

[Install]
WantedBy=multi-user.target
