# parent/systemd/mint-secrets-receiver.service
[Unit]
Description=Mint Enclave secrets receiver (first boot only)
After=mint-enclave.service
Requires=mint-enclave.service

[Service]
Type=oneshot
ExecStart=/opt/mint-enclave/receive-secrets.sh
StandardOutput=journal
StandardError=journal
TimeoutStartSec=180

[Install]
WantedBy=multi-user.target
