# parent/systemd/mint-secrets-sender.service
[Unit]
Description=Mint Enclave encrypted secrets sender (normal boot)
After=mint-creds-sender.service
Requires=mint-enclave.service

[Service]
Type=oneshot
# Script retries internally - enclave may not be listening on port 9002 yet
ExecStart=/opt/mint-enclave/send-secrets.sh
StandardOutput=journal
StandardError=journal
TimeoutStartSec=120

[Install]
WantedBy=multi-user.target
