# parent/systemd/mint-postgres-proxy.service
[Unit]
Description=Mint Enclave PostgreSQL vsock proxy
After=network.target postgresql.service nitro-enclaves-allocator.service
Wants=postgresql.service

[Service]
Type=simple
ExecStart=/usr/bin/vsock-proxy 5432 127.0.0.1 5432
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
